# addons/plex-now-showing/rootfs/etc/services.d/plex-now-showing/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Plex Now Showing add-on — service runner.
#
# Reads /data/options.json via bashio and forwards each option into the
# environment the Node server expects (see server/src/config.js for the
# authoritative list).
# ==============================================================================

set -e

bashio::log.info "Starting Plex Now Showing server..."

# ---- Home Assistant API (add-on mode) ----
# SUPERVISOR_TOKEN is injected by Supervisor when hassio_api / homeassistant_api
# are true in config.yaml — the server will auto-detect add-on mode.
export HA_URL="http://supervisor/core"
export HA_TOKEN="${SUPERVISOR_TOKEN}"

# Helper: export $1 from bashio config key $2. Assign then export (avoids
# masking the return value, which "export FOO=$(cmd)" does silently).
export_opt() {
  local _val
  _val="$(bashio::config "$2")"
  export "$1=${_val}"
}

# ---- User options → env ----
export_opt PLEX_URL          plex_url
export_opt PLEX_TOKEN        plex_token
export_opt PLEX_USERNAME     plex_username
export_opt PLEX_PLAYER       plex_player
export_opt LANDSCAPE         landscape
export_opt THEME             theme
export_opt POLL              poll_interval
export_opt STATE_TTL_MS      state_ttl_ms
export_opt MEDIA_INFO_TTL_MS media_info_ttl_ms
export_opt PROXY_SECRET      proxy_secret

# allowed_origins is a list → join with commas for the server's parser.
if bashio::config.has_value 'allowed_origins'; then
  ALLOWED_ORIGINS="$(bashio::config 'allowed_origins | join(",")')"
  export ALLOWED_ORIGINS
fi

# ---- Fully Kiosk auto-switcher (#48) ----
export_opt SWITCHER_ENABLED     switcher_enabled
export_opt SWITCHER_INTERVAL_MS switcher_interval_ms

# fully_kiosks is a structured list; fold into the server's pipe-delimited
# "host|password|playing_url[|stopped_url]" format, one kiosk per newline.
if bashio::config.has_value 'fully_kiosks'; then
  # jq — newline-joined so parseKiosks() happily splits on \n.
  FULLY_KIOSKS="$(bashio::config 'fully_kiosks | map([.host, .password, .playing_url, .stopped_url // ""] | map(select(. != "")) | join("|")) | join("\n")')"
  export FULLY_KIOSKS
fi

# ---- Diagnostics (no secrets) ----
bashio::log.info "Plex URL configured:      $([[ -n "${PLEX_URL}"   ]] && echo yes || echo no)"
bashio::log.info "Plex token configured:    $([[ -n "${PLEX_TOKEN}" ]] && echo yes || echo no)"
bashio::log.info "Pinned player:            ${PLEX_PLAYER:-<none, using username filter>}"
bashio::log.info "Plex username:            ${PLEX_USERNAME:-<none>}"
bashio::log.info "Landscape mode:           ${LANDSCAPE}"
bashio::log.info "Theme:                    ${THEME}"
bashio::log.info "Poll interval (ms):       ${POLL}"
bashio::log.info "State cache TTL (ms):     ${STATE_TTL_MS}"
bashio::log.info "Media-info cache TTL (ms):${MEDIA_INFO_TTL_MS}"
bashio::log.info "Proxy secret:             $([[ -n "${PROXY_SECRET}" ]] && echo set || echo unset)"
bashio::log.info "Allowed origins:          ${ALLOWED_ORIGINS:-<any (ingress-only)>}"
bashio::log.info "Switcher enabled:         ${SWITCHER_ENABLED:-false}"
if bashio::var.true "${SWITCHER_ENABLED:-false}"; then
  kiosk_count=$(bashio::config 'fully_kiosks | length')
  bashio::log.info "Switcher tick interval:   ${SWITCHER_INTERVAL_MS}ms"
  bashio::log.info "Kiosks configured:        ${kiosk_count}"
fi

cd /app
exec node src/server.js
