# addons/plex-now-showing/rootfs/etc/services.d/plex-now-showing/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Plex Now Showing add-on — service runner.
#
# Reads /data/options.json via bashio and forwards each option into the
# environment the Node server expects (see server/src/config.js for the
# authoritative list).
# ==============================================================================

set -e

bashio::log.info "Starting Plex Now Showing server..."

# ---- Home Assistant API (add-on mode) ----
# SUPERVISOR_TOKEN is injected by Supervisor when hassio_api / homeassistant_api
# are true in config.yaml — the server will auto-detect add-on mode.
export HA_URL="http://supervisor/core"
export HA_TOKEN="${SUPERVISOR_TOKEN}"

# Helper: export $1 from bashio config key $2. Assign then export (avoids
# masking the return value, which "export FOO=$(cmd)" does silently).
export_opt() {
  local _val
  _val="$(bashio::config "$2")"
  export "$1=${_val}"
}

# ---- User options → env ----
export_opt BACKEND           backend
export_opt DISPLAY_MODE      display_mode
export_opt PLAYER            player
export_opt PLEX_URL          plex_url
export_opt PLEX_TOKEN        plex_token
export_opt PLEX_USERNAME     plex_username
export_opt PLEX_PLAYER       plex_player
export_opt LANDSCAPE         landscape
export_opt THEME             theme
export_opt POLL              poll_interval
export_opt STATE_TTL_MS      state_ttl_ms
export_opt MEDIA_INFO_TTL_MS media_info_ttl_ms
export_opt COMING_SOON_TTL_MS coming_soon_ttl_ms
export_opt PROXY_SECRET      proxy_secret
export_opt COMING_SOON_TITLE coming_soon_title
export_opt RADARR_URL        radarr_url
export_opt RADARR_API_KEY    radarr_api_key
export_opt SONARR_URL        sonarr_url
export_opt SONARR_API_KEY    sonarr_api_key
export_opt COMING_SOON_MOVIES_COUNT coming_soon_movies_count
export_opt COMING_SOON_SHOWS_COUNT  coming_soon_shows_count
export_opt COMING_SOON_CYCLE_INTERVAL coming_soon_cycle_interval
export_opt COMING_SOON_DAYS_OFFSET coming_soon_days_offset
export_opt COMING_SOON_IMAGE_TYPE  coming_soon_image_type

# allowed_origins is a list → join with commas for the server's parser.
if bashio::config.has_value 'allowed_origins'; then
  ALLOWED_ORIGINS="$(bashio::config 'allowed_origins | join(",")')"
  export ALLOWED_ORIGINS
fi

# ---- Fully Kiosk auto-switcher (#48) ----
export_opt SWITCHER_ENABLED     switcher_enabled
export_opt SWITCHER_INTERVAL_MS switcher_interval_ms

# fully_kiosks is a structured list; fold into the server's pipe-delimited
# "host|password|playing_url[|stopped_url]" format, one kiosk per newline.
if bashio::config.has_value 'fully_kiosks'; then
  # jq — newline-joined so parseKiosks() happily splits on \n.
  FULLY_KIOSKS="$(bashio::config 'fully_kiosks | map([.host, .password, .playing_url, .stopped_url // ""] | map(select(. != "")) | join("|")) | join("\n")')"
  export FULLY_KIOSKS
fi

# ---- Visual toggles (V2 visual sprint) ----
export_opt VISUAL_PROGRESS_BAR        visual_progress_bar
export_opt VISUAL_RATINGS_BADGES      visual_ratings_badges
export_opt VISUAL_GENRE_CHIPS         visual_genre_chips
export_opt VISUAL_INFO_PANEL_MODE     visual_info_panel_mode
export_opt VISUAL_FRAME_STYLE         visual_frame_style
export_opt VISUAL_MARQUEE_FONT        visual_marquee_font
export_opt VISUAL_USE_BACKDROPS       visual_use_backdrops
export_opt VISUAL_BACKDROP_STYLE      visual_backdrop_style
export_opt VISUAL_BACKDROP_DELAY_MS   visual_backdrop_delay_ms
# Burn-in mitigation (#28). Master switch + sub-feature knobs.
export_opt VISUAL_BURN_IN_MITIGATION  visual_burn_in_mitigation
export_opt VISUAL_NUDGE_INTERVAL_MS   visual_nudge_interval_ms
export_opt VISUAL_NUDGE_AMPLITUDE_PX  visual_nudge_amplitude_px
export_opt VISUAL_NIGHT_MODE_ENTITY   visual_night_mode_entity
export_opt VISUAL_NIGHT_MODE_OPACITY  visual_night_mode_opacity
# Theme presets (#23) + accent colour picker (#66). Theme is one of
# classic-gold|art-deco-silver|neon-80s|minimalist-dark; accent colour is
# a strict #RRGGBB override (empty = use theme's default ramp).
export_opt VISUAL_THEME               visual_theme
export_opt VISUAL_ACCENT_COLOR        visual_accent_color
export_opt VISUAL_CORNER_RADIUS_PX    visual_corner_radius_px

# ---- Diagnostics (no secrets) ----
bashio::log.info "Backend:                  ${BACKEND:-plex}"
bashio::log.info "Display mode:             ${DISPLAY_MODE:-now_showing}"
bashio::log.info "Player entity:            ${PLAYER:-<auto>}"
bashio::log.info "Plex URL configured:      $([[ -n "${PLEX_URL}"   ]] && echo yes || echo no)"
bashio::log.info "Plex token configured:    $([[ -n "${PLEX_TOKEN}" ]] && echo yes || echo no)"
bashio::log.info "Legacy Plex player:       ${PLEX_PLAYER:-<none>}"
bashio::log.info "Plex username:            ${PLEX_USERNAME:-<none>}"
bashio::log.info "Landscape mode:           ${LANDSCAPE}"
bashio::log.info "Theme:                    ${THEME}"
bashio::log.info "Poll interval (ms):       ${POLL}"
bashio::log.info "State cache TTL (ms):     ${STATE_TTL_MS}"
bashio::log.info "Media-info cache TTL (ms):${MEDIA_INFO_TTL_MS}"
bashio::log.info "Coming Soon cache TTL:    ${COMING_SOON_TTL_MS:-900000}"
bashio::log.info "Coming Soon title:        ${COMING_SOON_TITLE:-Coming Soon}"
bashio::log.info "Radarr configured:        $([[ -n "${RADARR_URL}" && -n "${RADARR_API_KEY}" ]] && echo yes || echo no)"
bashio::log.info "Sonarr configured:        $([[ -n "${SONARR_URL}" && -n "${SONARR_API_KEY}" ]] && echo yes || echo no)"
bashio::log.info "Proxy secret:             $([[ -n "${PROXY_SECRET}" ]] && echo set || echo unset)"
bashio::log.info "Allowed origins:          ${ALLOWED_ORIGINS:-<any (ingress-only)>}"
bashio::log.info "Visual progress bar:      ${VISUAL_PROGRESS_BAR:-false}"
bashio::log.info "Visual ratings badges:    ${VISUAL_RATINGS_BADGES:-false}"
bashio::log.info "Visual genre chips:       ${VISUAL_GENRE_CHIPS:-false}"
bashio::log.info "Info panel mode:          ${VISUAL_INFO_PANEL_MODE:-on_tap}"
bashio::log.info "Frame style:              ${VISUAL_FRAME_STYLE:-bulbs}"
bashio::log.info "Marquee font:             ${VISUAL_MARQUEE_FONT:-bebas-neue}"
bashio::log.info "Backdrop art:             ${VISUAL_USE_BACKDROPS:-false} (${VISUAL_BACKDROP_STYLE:-fullscreen}, delay ${VISUAL_BACKDROP_DELAY_MS:-10000}ms)"
bashio::log.info "Burn-in mitigation:       ${VISUAL_BURN_IN_MITIGATION:-false} (nudge ${VISUAL_NUDGE_AMPLITUDE_PX:-4}px every ${VISUAL_NUDGE_INTERVAL_MS:-60000}ms, night entity=${VISUAL_NIGHT_MODE_ENTITY:-<none>}, opacity=${VISUAL_NIGHT_MODE_OPACITY:-0.4})"
bashio::log.info "Visual theme:             ${VISUAL_THEME:-classic-gold}"
bashio::log.info "Accent colour override:   ${VISUAL_ACCENT_COLOR:-<theme default>}"
bashio::log.info "Corner radius:            ${VISUAL_CORNER_RADIUS_PX:-0}px"
bashio::log.info "Switcher enabled:         ${SWITCHER_ENABLED:-false}"
if bashio::var.true "${SWITCHER_ENABLED:-false}"; then
  kiosk_count=$(bashio::config 'fully_kiosks | length')
  bashio::log.info "Switcher tick interval:   ${SWITCHER_INTERVAL_MS}ms"
  bashio::log.info "Kiosks configured:        ${kiosk_count}"
fi

cd /app
exec node src/server.js
